# package/base-files/files/etc/init.d/initscripts
#!/bin/sh /etc/rc.common

START=12
STOP=84

start()
{
	/sbin/eth_wdt &
	/etc/initscripts/symlinks &
	/etc/initscripts/modeminit &
	/etc/initscripts/usbhost &
	/etc/initscripts/gpsinit &
	#/etc/initscripts/upgrade &	
	/etc/initscripts/rtcinit
}
