# package/base-files/files/etc/hotplug.d/usb/02-sim5300_quirk
#!/bin/sh

if [ "$ACTION" == "add" -a "$PRODUCT" == "1519/20/1276" ]; then
    /etc/simman/setsim.sh -p  
fi

exit 0
