# package/base-files/files/etc/rc.local
ln -s /dev/ttyAPP4 /dev/com0
ln -s /dev/ttyAPP2 /dev/com1

/etc/initscripts/ethinit &

/etc/initscripts/beeper &

/etc/initscripts/update

exit 0